# etc/supervisord.conf
[program:server]
directory=/usr/local
command=/usr/local/bin/app
autostart=true
autorestart=true
stderr_logfile=/var/log/server.err
stdout_logfile=/var/log/server.log
environment=CODENATION_ENV=prod